# stage9/99-confinit/confinit/systemd/confinit-final@.service
[Unit]
Description=Confinit final files from /%I
DefaultDependencies=no
After=local-fs.target
Before=shutdown.target
Wants=confinit-intermediate@%i.service
After=confinit-intermediate@%i.service
Conflicts=shutdown.target
ConditionFileIsExecutable=/bin/confinit
ConditionPathExists=/%I/.confinit-final.yml
# RequiresMountsFor=/boot

[Service]
Type=oneshot
RemainAfterExit=yes
StandardOutput=journal
User=root
Group=root
ExecStart=/bin/confinit --config /%I/.confinit-final.yml
PrivateTmp=yes

[Install]
WantedBy=multi-user.target